# systemd/tg-dog.service
[Unit]
Description=tg-dog Telegram daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/tg-dog
Environment=APP_DIR=%h/tg-dog
ExecStart=%h/tg-dog/.venv/bin/tg-dog daemon --config %h/tg-dog/config.yaml --log-file %h/tg-dog/logs/daemon.log
ExecStartPre=/usr/bin/mkdir -p %h/tg-dog/logs %h/tg-dog/data %h/tg-dog/downloads %h/tg-dog/sessions
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
